# mpd_oled.service
[Unit]
Description=MPD OLED Display

[Service]
ExecStart=/usr/local/bin/mpd_oled -o 1 -b 24 -g 2 -f 30 -R

[Install]
WantedBy=multi-user.target
